# chapter5/chef-unreal-ircd/templates/default/init.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          ircd
# Required-Start:    networking $local_fs
# Required-Stop:     networking $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Unreal IRC daemon
# Description:       Unreal IRC server
### END INIT INFO


NAME=unreal-ircd
DESC="Unreal IRC Daemon"
DAEMON=<%= node['unreal']['source']['binary_destination'] %>
USER=<%= node['unreal']['user'] %>
PIDFILE="<%= node['unreal']['pidfile'] %>"

test -x $DAEMON || exit 0

set -e

case "$1" in
  start)
	echo -n "Starting $DESC: "
	if [ -f $PIDFILE ] && pgrep --pidfile $PIDFILE >/dev/null
	then
		echo "$NAME is already running ($(cat $PIDFILE))"
		exit 1
	else
		sudo -u $USER $DAEMON start
		echo "$NAME started with pid $(cat $PIDFILE)"
	fi
	echo "$NAME."
	;;
  status)
	if [ -f $PIDFILE ]
	then
		echo "$NAME is running ($(cat $PIDFILE))"
		exit 0
	else
		echo "$NAME is not running."
		exit 1
	fi
	;;
  stop)
	echo -n "Stopping $DESC: "

	sudo -u $USER $DAEMON stop > /dev/null || true
	pgrep --pidfile $PIDFILE && pkill --pidfile $PIDFILE
	rm -f $PIDFILE

	echo "$NAME."
	;;
 reload|force-reload|rehash)
	echo "Reloading $DESC configuration files."
	sudo -u $USER $DAEMON rehash
	;;
  restart)
	#
	#	If the "reload" option is implemented, move the "force-reload"
	#	option to the "reload" entry above. If not, "force-reload" is
	#	just the same as "restart".
	#
	echo -n "Restarting $DESC: "
	sudo -u $USER $DAEMON stop || true
	pgrep --pidfile $PIDFILE && pkill --pidfile $PIDFILE
	rm -f $PIDFILE

	sleep 1

	sudo -u $USER $DAEMON start

	echo "$NAME."
	;;
  *)
	N=/etc/init.d/$NAME
	# echo "Usage: $N {start|stop|restart|reload|rehash|force-reload}" >&2
	echo "Usage: $N {start|stop|restart|reload|rehash|force-reload}" >&2
	exit 1
	;;
esac

exit 0
